# assets/supervisord.conf
[supervisord]
nodaemon=true

[program:alfresco]
startsecs = 0
autorestart = false
command=/alfresco/init.sh
#user=alfresco
